# systemd/elevator-config-web.service.in
[Unit]
Description=Elevator RFID Config — local web admin tool (localhost:8080)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{PYTHON}} -m installer --web --no-browser --port 8080
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=elevator-config-web
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
