# systemd/fix_cp210x.service.in
[Unit]
Description=Fix cp210x RFID USB ID (bind 1560:0003 to the in-kernel driver)
After=multi-user.target

[Service]
Type=oneshot
# נשאר "active" אחרי שהסתיים, כדי ש-After=fix_cp210x.service של ה-tracker יסתמך עליו אמין.
RemainAfterExit=yes
ExecStart=/usr/local/bin/fix_cp210x.sh

[Install]
WantedBy=multi-user.target
