# systemd/rfid-tracker.service.in
[Unit]
Description=RFID Elevator Tracker - reads tags and reports floor to Firebase
After=network-online.target fix_cp210x.service
Wants=network-online.target fix_cp210x.service
# לא לוותר על הפעלה-מחדש לעולם: באתחול הפורט עשוי להופיע באיחור, וה-tracker
# אמור להמשיך לנסות עד שיצליח - במקום להיכנס ל-failed אחרי כמה ניסיונות.
StartLimitIntervalSec=0

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}
# מוודאים ש-rfid_config.json בבעלות {{USER}} לפני העלייה - אחרת עריכה/עדכון ידני
# כ-root משאיר אותו root-owned והשירות (רץ כ-{{USER}}) קורס ב-PermissionError.
# רק הקובץ עצמו (תחילית '+' = כ-root) - לא chown -R על כל הריפו/venv, שלקח כמה
# שניות על Pi Zero ובהן systemd מדווח 'activating' והדשבורד הבהב אדום בכל עלייה.
# בעלות .git ל-git pull עתידי מטופלת ב-setup.sh (chown אחרי pull).
ExecStartPre=+/bin/sh -c '[ -f {{PROJECT_DIR}}/rfid_config.json ] && chown {{USER}}:{{USER}} {{PROJECT_DIR}}/rfid_config.json 2>/dev/null || true'
# המתנה לקיום פורט הסיריאל (עד 60 שניות) לפני פתיחתו - מחליף את ה-sleep הקבוע,
# ומונע כשל באתחול כש-fix_cp210x עדיין מאתחל את ה-USB. ממשיכים גם אם לא הופיע,
# וה-tracker ייכשל-ויחזור (Restart=always) עד שהפורט יהיה זמין.
ExecStartPre=/bin/sh -c 'for i in $(seq 1 60); do [ -e {{SERIAL_PORT}} ] && exit 0; sleep 1; done; exit 0'
ExecStart={{PYTHON}} -u {{PROJECT_DIR}}/elevator_tracker_rfid.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rfid-tracker
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
