# systemd/shabbat-detector.service.in
[Unit]
Description=Shabbat Elevator Detector — behavioral Shabbat-mode detection
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}
# עריכה/עדכון ידני כ-root עלול להשאיר את הקונפיג בבעלות root (שובר את השירות
# ב-PermissionError) או את .git בבעלות root (שובר git pull עתידי של {{USER}}).
# תחילית '+' מריצה chown -R כ-root לכל הריפו לפני ההפעלה, מותנה בקיום setup.sh
# כבדיקת-שפיות שלא נוגעים בנתיב שגוי. כך הכל מתקן את עצמו בכל start, בלי setup.sh.
ExecStartPre=+/bin/sh -c '[ -f {{PROJECT_DIR}}/setup.sh ] && chown -R {{USER}}:{{USER}} {{PROJECT_DIR}} 2>/dev/null || true'
ExecStart={{PYTHON}} -m shabbat_detector.detector --config {{PROJECT_DIR}}/rfid_config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shabbat-detector
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
