# systemd/shabbat-detector.service.in
[Unit]
Description=Shabbat Elevator Detector — behavioral Shabbat-mode detection
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}
# הקונפיג עלול להישאר בבעלות root אחרי עריכה/עדכון ידני כ-root, ואז השירות (שרץ
# כ-{{USER}}) נכשל ב-PermissionError. תחילית '+' מריצה את ה-chown כ-root עוד לפני
# ההפעלה, כך שהשירות מתקן את עצמו בכל start - בלי תלות בהרצת setup.sh.
ExecStartPre=+/bin/sh -c 'chown {{USER}}:{{USER}} {{PROJECT_DIR}}/rfid_config.json 2>/dev/null || true'
ExecStart={{PYTHON}} -m shabbat_detector.detector --config {{PROJECT_DIR}}/rfid_config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shabbat-detector
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
