# systemd/shabbat-detector.service.in
[Unit]
Description=Shabbat Elevator Detector — behavioral Shabbat-mode detection
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_DIR}}
# מוודאים ש-rfid_config.json בבעלות {{USER}} לפני העלייה - אחרת עריכה/עדכון ידני
# כ-root משאיר אותו root-owned והשירות (רץ כ-{{USER}}) קורס ב-PermissionError.
# רק הקובץ עצמו (תחילית '+' = כ-root) - לא chown -R על כל הריפו/venv, שלקח כמה
# שניות על Pi Zero ובהן systemd מדווח 'activating' והדשבורד הבהב אדום בכל עלייה.
# בעלות .git ל-git pull עתידי מטופלת ב-setup.sh (chown אחרי pull).
ExecStartPre=+/bin/sh -c '[ -f {{PROJECT_DIR}}/rfid_config.json ] && chown {{USER}}:{{USER}} {{PROJECT_DIR}}/rfid_config.json 2>/dev/null || true'
ExecStart={{PYTHON}} -m shabbat_detector.detector --config {{PROJECT_DIR}}/rfid_config.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shabbat-detector
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
